# deploy/systemd/gmgn-twitter-cli.service
[Unit]
Description=GMGN Twitter CLI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=qinghuan
Group=staff
WorkingDirectory=/Users/qinghuan/Documents/code/GmgnTwitterTgAlert
EnvironmentFile=/Users/qinghuan/Documents/code/GmgnTwitterTgAlert/.env
ExecStart=/Users/qinghuan/Documents/code/GmgnTwitterTgAlert/.venv/bin/gmgn-twitter-cli serve
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/Users/qinghuan/Documents/code/GmgnTwitterTgAlert/data /Users/qinghuan/Documents/code/GmgnTwitterTgAlert/logs

[Install]
WantedBy=multi-user.target
